# systemd/tailscale-node-online@.target
[Unit]
Description=Synced with tailscale status

BindsTo=tailscaled.service
After=tailscaled.service
